# ansible/roles/common/templates/node-exporter.service.j2
[Unit]
Description=Prometheus node_exporter
Documentation=https://prometheus.io/docs/guides/node-exporter/
After=network.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/node_exporter \
  --collector.systemd \
  --collector.processes \
  --web.listen-address=:9100 \
  --web.telemetry-path=/metrics

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=full
ProtectHome=yes
# CapabilityBoundingSet must include CAP_SYS_PTRACE for --collector.processes
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

Restart=on-failure
RestartSec=5s
SendSIGKILL=no
TimeoutStopSec=20s

[Install]
WantedBy=multi-user.target
